# mcp-bridge.service
[Unit]
Description=MCP Bridge Service
After=network.target

[Service]
Type=simple
ExecStart=/home/karld@tusk-watf.tuskerdirect.com/Nextcloud/Tusker/code/mcp/mcp-bridge/mcp-bridge
WorkingDirectory=/home/karld@tusk-watf.tuskerdirect.com/Nextcloud/Tusker/code/mcp/mcp-bridge
EnvironmentFile=/home/karld@tusk-watf.tuskerdirect.com/Nextcloud/Tusker/code/mcp/mcp-bridge/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target